# scripts/ssh-agent.service
[Unit]
Description=ssh agent

[Service]
ExecStart=/usr/bin/ssh-agent -D

[Install]
WantedBy=default.taraget
